# backend/deploy/aimeditation.service
[Unit]
Description=AI Meditation Backend Service
After=network.target

[Service]
Type=simple
User=aimeditation
Group=aimeditation
WorkingDirectory=/var/www/aimeditation/backend
Environment=PATH=/var/www/aimeditation/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/var/www/aimeditation/backend
Environment=APP_ENV=testing

ExecStart=/var/www/aimeditation/backend/venv/bin/gunicorn main:app \
    --bind 0.0.0.0:8000 \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --access-logfile /var/log/aimeditation/access.log \
    --error-logfile /var/log/aimeditation/error.log \
    --capture-output \
    --enable-stdio-inheritance

Restart=always
RestartSec=5
StartLimitInterval=0

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
CapabilityBoundingSet=
AmbientCapabilities=
SecureBits=no-setuid-fixup-locked

[Install]
WantedBy=multi-user.target 